# deploy_server.crontab
MAILTO=""
*/5 * * * * /home/deploy_server/venv/bin/python /home/deploy_server/script/shop_generator.py && /home/deploy_server/venv/bin/python /home/deploy_server/script/upload_to_database.py
